# rootconf/etc/profile.d/opx.sh
if [ -r /etc/opx/opx-environment.sh ]
then
    source /etc/opx/opx-environment.sh
fi

